# image/files/systemd/trailcurrent-playbill-firstboot.service
[Unit]
Description=TrailCurrent Playbill first-boot setup (rootfs expand, machine-id, SSH host keys)
After=local-fs.target
Before=ssh.service
ConditionPathExists=!/var/lib/trailcurrent-playbill/.firstboot-done

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/trailcurrent-playbill-firstboot.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
